# weather-report.timer
[Unit]
Description=Emma's Daily Weather Briefing Timer

[Timer]
OnCalendar=*-*-* 06:45:00
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target